# supervisord.conf


[supervisord]
logfile=/workspaces/dominion-os-demo-build/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/workspaces/dominion-os-demo-build/logs/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=vscode

[program:oauth_server]
command=/workspaces/dominion-os-demo-build/.venv/bin/gunicorn -w 4 --timeout 60 --log-level info -b 127.0.0.1:5000 app:app
directory=/workspaces/dominion-os-demo-build/oauth_server
autostart=true
autorestart=true
stderr_logfile=/workspaces/dominion-os-demo-build/oauth_server/gunicorn.err.log
stdout_logfile=/workspaces/dominion-os-demo-build/oauth_server/gunicorn.out.log
user=vscode
healthcheck_program=/workspaces/dominion-os-demo-build/.venv/bin/python /workspaces/dominion-os-demo-build/oauth_server/health_check.py
healthcheck_interval=30

[program:widget_service]
command=/workspaces/dominion-os-demo-build/.venv/bin/gunicorn -w 4 --timeout 60 --log-level info -b 127.0.0.1:8080 app:app
directory=/workspaces/dominion-os-demo-build/widget_service
autostart=true
autorestart=true
stderr_logfile=/workspaces/dominion-os-demo-build/widget_service/gunicorn.err.log
stdout_logfile=/workspaces/dominion-os-demo-build/widget_service/gunicorn.out.log
user=vscode
healthcheck_program=/workspaces/dominion-os-demo-build/.venv/bin/python /workspaces/dominion-os-demo-build/widget_service/health_check.py
healthcheck_interval=30
